# config/chroot_local-includes/lib/systemd/system/gobby-server.service
[Unit]
Description=Gobby Server Daemon
After=network.target

[Service]
ExecStart=/usr/bin/infinoted
